# crontab/global.properties
#a.name.last=jie
#b.name.first=shi
#b.name=shijie

*/5 * * * *=dir c:
#*/2 * * * *=dir d: